# stonks-watcher/stonks-scraper/crontab
PATH=/usr/local/bin

*/2 * * * * root cd /app && scrapy crawl olx_spider >> /var/log/cron.log 2>&1
# Don't remove the empty line at the end of this file. It is required to run the cron job